# osh/hub/osh-retention.service
[Unit]
Description=OpenScanHub retention policy enforcer

[Service]
ExecStart=/usr/sbin/osh-retention
Type=oneshot

[Install]
WantedBy=multi-user.target
