# osh/hub/osh-stats.service
[Unit]
Description=OpenScanHub statistics collector

[Service]
ExecStart=/usr/bin/osh-stats
Type=oneshot

[Install]
WantedBy=multi-user.target
